# packages/o/openssh.files/sshd-keygen.service
[Unit]
Description=OpenSSH Server Key Generation
ConditionPathExists=|!/etc/ssh/ssh_host_rsa_key

[Service]
Environment="SSH_USE_STRONG_RNG=14"
ExecStart=/bin/ssh-keygen -A
Type=oneshot
