# packages/i/iptables.files/iptables.service
[Unit]
Description=iptables firewall
Before=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/etc/iptables.sh

[Install]
WantedBy=network.target
